# roles/pve/wake-on-lan/templates/wol.service.j2
[Unit]
Description=Enable Wake On Lan

[Service]
Type=oneshot
ExecStart = /usr/sbin/ethtool --change {{ wol_network_interface_name }} wol g

[Install]
WantedBy=basic.target